# 04_keyboard/99-hello-keyboard.rules
SUBSYSTEM=="usb", ACTION=="add", ATTRS{bInterfaceClass}=="03", ATTRS{bInterfaceProtocol}=="01" RUN+="/sbin/modprobe hello_keyboard"
